# config/systemd/dictation-service.service
[Unit]
Description=GPU-Optimized Dictation Service
After=graphical-session.target

[Service]
Type=simple
Environment="PATH=%h/miniconda3/envs/whisper/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CONDA_PREFIX=%h/miniconda3/envs/whisper"
Environment="DISPLAY=:0"
# Try multiple possible Xauthority locations
Environment="XAUTHORITY=/run/user/%U/gdm/Xauthority"
ExecStartPre=/bin/bash -c 'export DISPLAY=:0; export XAUTHORITY=${XAUTHORITY:-%h/.Xauthority}; xhost +local: > /dev/null 2>&1 || true'
# Use wrapper script that handles X11/Wayland properly
ExecStart=/bin/bash %h/.local/share/dictation-service/dictation-wrapper.sh
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.local/share/dictation-service/logs/dictation-service.log
StandardError=append:%h/.local/share/dictation-service/logs/dictation-service.log

[Install]
WantedBy=default.target